# rng.service
[Unit]
Description=Send random numbers generated by OneRNG through WebSockets
After=network.target

[Service]
ExecStartPre=/bin/bash -c "stty raw -echo -F /dev/ttyACM0; echo cmd0 > /dev/ttyACM0; echo cmdO > /dev/ttyACM0"
ExecStart=/usr/local/bin/rng
Restart=always
User=root
Group=root
Environment="PATH=/usr/bin:/usr/local/bin:/bin"
Environment="NODE_ENV=production"
Environment="ONERNG_PATH=/dev/ttyACM0"

[Install]
WantedBy=multi-user.target
