# rng.service
[Unit]
Description=Send random number generated by OneRNG throw WebSockets

[Service]
ExecStart=/usr/bin/rng
Restart=always
User=root
Group=root
Environment=PATH=/usr/bin:/usr/local/bin
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target