# rng.service
[Unit]
Description=Send random number generated by OneRNG through WebSockets

[Service]
ExecStart=/bin/bash -c "stty raw -echo </dev/ttyACM0;echo cmd0 >/dev/ttyACM0;echo cmdO >/dev/ttyACM0;/usr/local/bin/rng"
Restart=always
User=root
Group=root
Environment=PATH=/usr/bin:/usr/local/bin
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
